# agent/box-agent.service
[Unit]
Description=bux box-agent (cloud control + agent task dispatch)
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/.config/bux/env

[Service]
Type=simple
EnvironmentFile=%h/.config/bux/env
WorkingDirectory=%h/.local/share/bux/repo/agent
ExecStart=%h/.local/share/bux/venv/bin/python %h/.local/share/bux/repo/agent/box_agent.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
