# agent/bux-tg.service
[Unit]
Description=bux Telegram bot (user-owned)
After=network-online.target
ConditionPathExists=%h/.config/bux/tg.env

[Service]
Type=simple
EnvironmentFile=%h/.config/bux/tg.env
WorkingDirectory=%h/.local/share/bux/repo/agent
ExecStart=%h/.local/share/bux/venv/bin/python %h/.local/share/bux/repo/agent/telegram_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
