# agent/bux-browser-keeper.service
[Unit]
Description=bux browser-keeper (local Chrome CDP browser)
After=network-online.target box-agent.service
Wants=network-online.target
ConditionPathExists=%h/.config/bux/env

[Service]
Type=simple
EnvironmentFile=%h/.config/bux/env
WorkingDirectory=%h/.local/share/bux/repo/agent
ExecStart=%h/.local/share/bux/venv/bin/python %h/.local/share/bux/repo/agent/browser_keeper.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
